# ansible/roles/technitium_lxc/files/technitium-zone-sync.timer
[Unit]
Description=Technitium DNS zone sync timer (every 15 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Unit=technitium-zone-sync.service

[Install]
WantedBy=timers.target
