# ansible/roles/technitium_lxc/files/technitium-zone-sync.service
[Unit]
Description=Technitium DNS zone sync from Git
After=network-online.target dns.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-zone.sh
StandardOutput=journal
StandardError=journal
# Fail the unit visibly so the timer shows a failed state
