# unraid/flash-config/plugins/dynamix/language-check.cron
# Generated languages version check schedule:
10 0 1 * * /usr/local/emhttp/plugins/dynamix.plugin.manager/scripts/languagecheck &> /dev/null

