# unraid/flash-config/plugins/dynamix/mover.cron
# Generated mover schedule:
0 */2 * * * /usr/local/sbin/mover start > /dev/null 2> >(logger -t move)

