# unraid/flash-config/plugins/dynamix/plugin-check.cron
# Generated plugins version check schedule:
10 */6 * * * /usr/local/emhttp/plugins/dynamix.plugin.manager/scripts/plugincheck &> /dev/null

